# roles/juventas/files/10-trackpoint.rules
SUBSYSTEM=="input", ATTR{name}=="*TrackPoint*", RUN+="/etc/rc.local"
