# roles/ananke/files/trackpoint.service
[Unit]
Description=Sets custom TrackPoint settings

[Service]
ExecStart=/opt/configure-trackpoint.sh

[Install]
WantedBy=multi-user.target
